# templates/exhibitor.service.erb
[Unit]
Description=Exhibitor service

[Service]
Type=oneshot
ExecStart=/opt/exhibitor/exhibitor.sh
StandardOutput=journal
ExecStop=kill %1

[Install]
