# templates/exhibitor.service.erb
[Unit]
Description=Exhibitor service

[Service]
Type=oneshot
ExecStart=/opt/exhibitor/exhibitor.sh start
StandardOutput=journal
ExecStop=/opt/exhibitor/exhibitor.sh stop

[Install]
